# data/etc/init/icms-services.conf
description "Services for ImpressCMS"
author "Raimondas RimkeviÄius (aka MekDrop)"

start on (net-device-up and started mysql and started cron and started ssh)

script
 sudo /usr/local/bin/icms-extract-backup.sh
end script